# recipes-core/switch-log-location/files/switch-log-location.service
[Unit]
Description=Setup log location (persistent or volatile)
DefaultDependencies=no
Before=systemd-journald.service
After=local-fs.target
RequiresMountsFor=/mnt/datastore

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/switch-log-location.sh

[Install]
WantedBy=local-fs.target
