# recipes-core/save-rtc/save-rtc/save-timestamp.service
[Unit]
Description=Save system clock to timestamp file

[Service]
Type=oneshot
ExecStart=/usr/bin/save-timestamp.sh
